# config/scheduling/claude-loop.service
# systemd service unit for Claude Code autonomous loop
#
# Install:
#   sudo cp config/scheduling/claude-loop.service /etc/systemd/system/
#   sudo cp config/scheduling/claude-loop.timer   /etc/systemd/system/
#   sudo systemctl enable --now claude-loop.timer
#
# Status:
#   systemctl status claude-loop.timer
#   journalctl -u claude-loop.service -f

[Unit]
Description=Claude Code Autonomous Roadmap Loop
After=network.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/path/to/project
ExecStart=/path/to/project/scripts/run.sh --all
EnvironmentFile=/path/to/project/.env
StandardOutput=append:/path/to/project/.claude-loop/systemd.log
StandardError=append:/path/to/project/.claude-loop/systemd.log
TimeoutStartSec=3600

---

# /etc/systemd/system/claude-loop.timer
[Unit]
Description=Run Claude Code loop every 2 hours

[Timer]
OnBootSec=5min
OnUnitActiveSec=2h
Persistent=true

[Install]
WantedBy=timers.target
